# etc/service/loopy/run
#!/bin/sh

exec /loopy.py /var/run/loopy.pid | logger -t loopy
